# supervisord_mysqld.conf
[program:mysqld]
command=/opt/start_mysqld.sh
numprocs=1
autostart=true
autorestart=unexpected
