# supervisord_mysqld.conf
[program:mysqld]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /etc/init.d/mysql start
numprocs=1
autostart=true
autorestart=true
