# supervisord_apache2.conf
[program:apache2]
command=/opt/start_apache2.sh
numprocs=1
autostart=true
autorestart=true
